# task/crontab
## update to match your tizmezone
TZ=America/New_York
0,15,30,45 * * * * cd /app && /usr/local/bin/python maker-world-api.py